# rootfs/etc/s6-overlay/s6-rc.d/svc-forgejo-runner/finish
#!/usr/bin/execlineb -S0

echo "$0" > /run/s6-linux-init-container-results/exitcode
exec /run/s6/basedir/bin/halt