# root/etc/periodic/weekly/docker-system-prune.sh
#!/bin/sh

exec docker system prune --all --force
