# root/etc/periodic/weekly/docker-system-prune.sh
#!/bin/sh

exec docker image prune --all --force
