# root/etc/periodic/daily/docker-volume-prune.sh
#!/bin/sh

exec docker volume prune --all --force
