# root/etc/s6-overlay/s6-rc.d/svc-forgejo-runner/finish
#!/command/with-contenv sh
# shellcheck shell=sh

exec /run/s6/basedir/bin/halt
