# rootfs/etc/s6-overlay/s6-rc.d/svc-docker/finish
#!/bin/sh

if test "$1" -eq 256 ; then
  e=$((128 + $2))
else
  e="$1"
fi

echo "$1" > /run/s6-linux-init-container-results/exitcode