# rootfs/etc/s6-overlay/s6-rc.d/svc-docker/run
#!/usr/bin/env bash

# Enable subreaping, because dind is running within s6
export TINI_SUBREAPER=1

exec /usr/local/bin/dockerd-entrypoint.sh
