# root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/command/with-contenv sh
# shellcheck shell=sh

exec crond -f -l 4
